# config/com.lean-bubbles.image-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lean-bubbles.image-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            if command -v bubble &amp;&gt; /dev/null; then
                bubble images build lean-base 2&gt;&amp;1 | logger -t lean-bubbles-images
            fi
        </string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/lean-bubbles-image-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/lean-bubbles-image-refresh.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
