# config/com.bubble.relay-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bubble.relay-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            if command -v bubble &amp;&gt; /dev/null; then
                bubble relay daemon 2&gt;&amp;1 | logger -t bubble-relay
            fi
        </string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/bubble-relay-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/bubble-relay-daemon.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
